# Rinputer2.service
[Unit]
Description=Ragnarok Input Daemon

[Service]
ExecStart=/usr/bin/rinputer2

[Install]
WantedBy=multi-user.target
